# openwrt/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
tty1::askfirst:/bin/ash --login
tty0::askfirst:/bin/ash --login
console::askfirst:/bin/ash --login
